# openrc/avahi2dns
#!/sbin/openrc-run

output_logger=logger
error_logger=logger
command_background=true
command="avahi2dns"
pidfile="/run/avahi2dns.pid"
