# openrc/avahi2dns
#!/sbin/openrc-run

output_logger="logger -t avahi2dns -p daemon.info"
error_logger="logger -t avahi2dns -p daemon.info"
command_background=true
command="avahi2dns"
pidfile="/run/avahi2dns.pid"
